# crontab
0 8 * * * cd /code && flask --app=webapp resend-notifications
