# crontab
0 8 * * * cd /code && uv run flask --app=webapp resend-notifications
