# crontab
0 8 * * * cd /code && flask resend-notifications
